# nas_ports/freenas/freenas-etcd/pkg-plist
/usr/local/etc/rc.d/etcd
/usr/local/sbin/etcd
