# nas_ports/freenas/freenas-dispatcher/pkg-plist
lib/dispatcher/src/frontend/__init__.py
lib/dispatcher/src/frontend/frontend.py
lib/dispatcher/src/frontend/templates/nav.html
lib/dispatcher/src/frontend/templates/tasks.html
lib/dispatcher/src/frontend/templates/events.html
lib/dispatcher/src/frontend/templates/rpc.html
lib/dispatcher/src/frontend/templates/term.html
lib/dispatcher/src/frontend/templates/stats.html
lib/dispatcher/src/frontend/templates/syslog.html
lib/dispatcher/src/frontend/templates/vm.html
lib/dispatcher/src/frontend/templates/filebrowser.html
lib/dispatcher/src/frontend/templates/base.html
lib/dispatcher/src/frontend/templates/apidoc/index.html
lib/dispatcher/src/frontend/templates/apidoc/nav.html
lib/dispatcher/src/frontend/templates/apidoc/rpc.html
lib/dispatcher/src/frontend/templates/apidoc/events.html
lib/dispatcher/src/frontend/templates/apidoc/tasks.html
lib/dispatcher/src/frontend/templates/apidoc/schemas.html
lib/dispatcher/src/frontend/static/middleware.js
lib/dispatcher/src/frontend/static/contextmenu.js
lib/dispatcher/src/frontend/static/term.js
lib/dispatcher/src/frontend/static/jquery.jsonview.js
lib/dispatcher/src/frontend/static/jquery.jsonview.css
lib/dispatcher/src/frontend/static/main.css
lib/dispatcher/src/frontend/static/filebrowser.css
lib/dispatcher/src/frontend/static/css/dashboard.css
lib/dispatcher/src/frontend/static/css/sidebar.css
lib/dispatcher/src/frontend/static/css/messenger-theme-future.css
lib/dispatcher/src/frontend/static/css/messenger.css
lib/dispatcher/src/frontend/static/js/angular-messenger.js
lib/dispatcher/src/frontend/static/js/app.js
lib/dispatcher/src/frontend/static/js/controllers.js
lib/dispatcher/src/frontend/static/js/directives.js
lib/dispatcher/src/frontend/static/js/messenger-theme-future.js
lib/dispatcher/src/frontend/static/js/messenger.min.js
lib/dispatcher/src/frontend/static/lib/WebsocketClient.js
lib/dispatcher/src/frontend/static/partials/apidocs.html
lib/dispatcher/src/frontend/static/partials/events.html
lib/dispatcher/src/frontend/static/partials/index.html
lib/dispatcher/src/frontend/static/partials/rpc.html
lib/dispatcher/src/frontend/static/partials/syslog.html
lib/dispatcher/src/frontend/static/partials/tests.html
lib/dispatcher/src/frontend/static/partials/stats.html
libexec/dispatcher/notify-pam
libexec/dispatcher/ssh_transport_catcher
lib/dispatcher/src/lib/sendzfs.so
etc/rc.d/dispatcher
sbin/dispatcherd
sbin/dispatcherctl
sbin/sdsinit
libexec/taskworker
etc/middleware.conf
