# nas_ports/freenas/freenas-dispatcher/pkg-plist
libexec/dispatcher/ssh_transport_catcher
etc/rc.d/dispatcher
sbin/dispatcherd
sbin/dispatcherctl
sbin/sdsinit
libexec/taskworker
etc/middleware.conf
