# nas_ports/freenas/freenas-dscached/pkg-plist
/usr/local/etc/rc.d/dscached
/usr/local/sbin/dscached
/usr/local/sbin/dscachectl
