# nas_ports/freenas/freenas-debugd/pkg-plist
/usr/local/etc/rc.d/debugctl.conf
/usr/local/etc/rc.d/debugd
/usr/local/sbin/debugd
/usr/local/sbin/debugctl
