# nas_ports/freenas/freenas-alertd/pkg-plist
/usr/local/etc/rc.d/alertd
/usr/local/sbin/alertd
