# nas_ports/freenas/freenas-containerd/pkg-plist
etc/rc.d/containerd
sbin/containerd
share/containerd/rancheros/vmlinuz
share/containerd/rancheros/initrd
