# nas_ports/freenas/freenas-containerd/pkg-plist
etc/rc.d/containerd
sbin/containerd
share/containerd/firmware/BHYVE_UEFI_20160526.fd
